# .github/docker/ptclient/etc-supervisor-conf.d-services.conf
[program:node]
directory = /gt/sc-prog-repos/emr/ptclient
command = /bin/bash -c "npm install && npm run start-nuxt-in-test-env"
user = root
autostart = true