# .github/docker-used-by-workflow/ptserver/etc-supervisor-conf.d-services.conf
[program:logs-permission]
command = chmod -R 777 /lumen-for-api/storage/
user = root
autostart = true
autorestart = false

[program:bootstrap-cache-permission]
command = chmod -R 777 /lumen-for-api/bootstrap/cache
user = root
autostart = true
autorestart = false

[program:external-packages]
directory = /lumen-for-api
command = composer install
user = root
autostart = true
autorestart = false

[program:ptserver]
directory = /lumen-for-api
command = /bin/bash -c "php -S localhost:8000 -t ."
user = root
autostart = true
