# .github/docker-used-by-workflow/ptserver/etc-supervisor-conf.d-services.conf
[program:logs-permission]
command = chmod -R 777 /gt/sc-prog-repos/emr/servers/lumen-for-api/storage/
user = root
autostart = true
autorestart = false

[program:bootstrap-cache-permission]
command = chmod -R 777 /gt/sc-prog-repos/emr/servers/lumen-for-api/bootstrap/cache
user = root
autostart = true
autorestart = false

[program:external-packages]
directory = /gt/sc-prog-repos/emr/servers/lumen-for-api
command = composer update
user = root
autostart = true
autorestart = false


[program:ptserver]
directory = /gt/sc-prog-repos/emr/servers
command = /bin/bash -c "npm install && npm run start-test"
user = root
autostart = true
