# .github/docker-used-by-workflow/puppet/etc-supervisor-conf.d-services.conf
[program:puppet]
directory = /gt/sc-prog-repos/emr/puppet
command = /bin/bash 
user = root
autostart = true
autorestart=false


